# prometheus-1w.conf
[supervisord]
[program:prometheus-1w]
command = /home/jgross/dev/prometheus-1w/prometheus-1w.py -v
stdout_logfile = /var/log/supervisor/%(program_name)s.log
stderr_logfile = /var/log/supervisor/%(program_name)s-stderr.log
autorestart = true
